# configs/docker/wireshark.conf
[program:wireshark]
priority=1
environment=DISPLAY=:1
command=/usr/bin/wireshark
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

